# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; Wright API (FastAPI + uvicorn)
[program:wright-api]
command=/usr/local/bin/uv run --project /workspace uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/workspace
environment=HOME="/home/agent",HERMES_HOME="/home/agent/.hermes",HERMES_PROFILE="wright",HERMES_API_BASE_URL="http://127.0.0.1:8642",HERMES_API_KEY="wright-local-dev-key-000000000000000000000000",API_SERVER_KEY="wright-local-dev-key-000000000000000000000000",PATH="/opt/hermes/bin:/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/wright-api-stdout.log
stderr_logfile=/var/log/supervisor/wright-api-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; Hermes Gateway (agent API on port 8642)
[program:hermes-gateway]
command=/opt/hermes/bin/hermes -p wright gateway run
environment=HOME="/home/agent",HERMES_HOME="/home/agent/.hermes",HERMES_PROFILE="wright",API_SERVER_KEY="wright-local-dev-key-000000000000000000000000",PATH="/opt/hermes/bin:/opt/hermes/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/hermes-gateway-stdout.log
stderr_logfile=/var/log/supervisor/hermes-gateway-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

